# docker/prod/services/php-cli/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:laravel-queue]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php -d variables_order=EGPCS /var/www/artisan queue:work --tries=3 --verbose
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log
stderr_logfile=/var/www/storage/logs/worker-errors.log
stopwaitsecs=3600

[program:telegram-queue]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php -d variables_order=EGPCS /var/www/artisan queue:work telegram --tries=1 --verbose
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/telegram-worker.log
stderr_logfile=/var/www/storage/logs/telegram-worker-errors.log
stopwaitsecs=3600

[program:notification-queue]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php -d variables_order=EGPCS /var/www/artisan queue:work notification --tries=3 --verbose
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/notification-worker.log
stderr_logfile=/var/www/storage/logs/notification-worker-errors.log
stopwaitsecs=3600

[program:laravel-scheduler]
process_name=%(program_name)s_%(process_num)02d
command=/bin/sh -c "while [ true ]; do (/usr/bin/php -d variables_order=EGPCS /var/www/artisan schedule:run --verbose --no-interaction &); sleep 60; done"
autostart=true
autorestart=true
stdout_logfile_maxbytes=0
stderr_logfile=/var/www/storage/logs/scheduler-errors.log
stderr_logfile_maxbytes=0
stdout_logfile=/var/www/storage/logs/scheduler.log
stopwaitsecs=3600

[program:telegram]
command=/usr/bin/php -d variables_order=EGPCS /var/www/html/artisan telegram:webhook --setup
autostart=true
autorestart=true
startretries=3
